# meta-ec-core/recipes-core/base-files/base-files/etc/profile.d/path.sh
#!/bin/sh

export PATH=@EC_TARGET_INSTALL_DIR@/bin:~/.local/bin:~/.diff-so-fancy:$PATH
